# devel/icecream/files/icecc-scheduler.in
#!/bin/sh
#
# $FreeBSD$

# PROVIDE: icecc_scheduler
# REQUIRE: NETWORKING ldconfig
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable icecc_scheduler:
#
# icecc_scheduler_enable (bool):        Set to NO to by default.
# icecc_scheduler_netname (str):        Netname. No default value.
# icecc_scheduler_port (int):           Port. Default value: 
# icecc_scheduler_logfile (str):        Log file. Default value: /var/log/icecc_scheduler.log
# icecc_scheduler_persistent (bool):    Persistent client connection. Default value: NO

. /etc/rc.subr

name=icecc_scheduler
rcvar=icecc_scheduler_enable

load_rc_config $name

: ${icecc_scheduler_enable:=NO}
: ${icecc_scheduler_netname:=}
: ${icecc_scheduler_port:=}
: ${icecc_scheduler_logfile:="/var/log/icecc_scheduler.log"}
: ${icecc_scheduler_persistent:=NO}

icecc_scheduler_enable=${icecc_scheduler_enable-"NO"}
icecc_scheduler_flags=${icecc_scheduler_flags-"-u icecc -d -l ${icecc_scheduler_logfile}"}
if [ -n "$icecc_scheduler_netname" ]; then
    icecc_scheduler_flags=${icecc_scheduler_flags}" -n ${icecc_scheduler_netname}"
fi
if [ -n "$icecc_scheduler_port" ]; then
    icecc_scheduler_flags=${icecc_scheduler_flags}" -p ${icecc_scheduler_port}"
fi
if checkyesno icecc_scheduler_persistent; then
    icecc_scheduler_flags=${icecc_scheduler_flags}" -r"
fi
command=%%PREFIX%%/sbin/icecc-scheduler
pidfile=${pidfile:-%%ICECC_SCHEDULER_PIDFILE%%}

icecc_scheduler_precmd()
{
	touch ${pidfile}
	chown icecc:icecc ${pidfile}
}

start_precmd="icecc_scheduler_precmd"
run_rc_command "$1"
